# drafts/report-platform/ai.eugot.zq-reports.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  DRAFT — inspect before installing. launchd LaunchAgent for the ZQ reports
  platform on himalayas (Evgeniis-Mac-mini, macOS).

  WHY launchd, not nohup: the prior `nohup ... &` server is orphaned (PPID=1)
  and does NOT survive logout/reboot. A LaunchAgent with KeepAlive restarts
  the server on crash and at login, giving durable hosting. (himalayas recon:
  the existing :4000 ziglint_dash is exactly such an orphan — this is the fix.)

  Install (manual, after inspection):
    cp ai.eugot.zq-reports.plist ~/Library/LaunchAgents/
    launchctl load -w ~/Library/LaunchAgents/ai.eugot.zq-reports.plist
  Status / logs:
    launchctl list | grep zq-reports
    tail -f /Volumes/Crucial/Data/zq-reports/server.log
  Uninstall:
    launchctl unload -w ~/Library/LaunchAgents/ai.eugot.zq-reports.plist

  NOTES:
    - Port 4010 (NOT 4000 — taken by ziglint_dash; NOT 8765 — Logitech).
    - Logs to /Volumes/Crucial (internal disk is 93% full per recon).
    - elixir pinned to 1.20.1-otp-29 to avoid the host's two-mise-root split.
    - mise path is the nix-darwin one: /etc/profiles/per-user/etretiakov/bin/mise
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.eugot.zq-reports</string>

  <key>ProgramArguments</key>
  <array>
    <string>/etc/profiles/per-user/etretiakov/bin/mise</string>
    <string>x</string>
    <string>elixir@1.20.1-otp-29</string>
    <string>--</string>
    <string>elixir</string>
    <string>/Users/etretiakov/zq-report/report-server.exs</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/etretiakov/zq-report</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>ZQ_REPORT_PORT</key>
    <string>4010</string>
    <key>ZQ_REPORT_HOST</key>
    <string>100.100.39.44</string>
    <key>ZQ_REPORTS_DIR</key>
    <string>/Users/etretiakov/zq-report/reports</string>
    <key>MISE_TRUSTED_CONFIG_PATHS</key>
    <string>/Users/etretiakov/zq-report</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>
  <key>ThrottleInterval</key>
  <integer>10</integer>

  <!-- Logs go to the HOME dir, not the external volume: a LaunchAgent's
       StandardErrorPath must be writable at spawn time, and an external
       /Volumes mount may be unavailable/again-permissioned in launchd's
       context → the job fails with exit 78 (EX_CONFIG) and writes nothing.
       Report DATA still lives on /Volumes/Crucial (ZQ_REPORTS_DIR). -->
  <key>StandardOutPath</key>
  <string>/Users/etretiakov/zq-report/server.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/etretiakov/zq-report/server.err.log</string>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
